# config-pkg.conf
description "first boot/config of the couchbase server AMI"
author      "couchbase"

start on startup
console output
exec /home/ec2-user/config-pkg

